# shutdown-hook.conf
# put this script under /etc/init/

description "run at shutdown"
start on starting rc
task

# Below is the path for the cron/shell script that you created.
exec /bin/bash preserve_logs.sh
